# roles/app/app/files/supervisord-gearman.ini
[program:gearman-statink]
command=/usr/bin/scl enable php82 -- /var/www/statink/statink/yii gearman/start %(process_num)s
process_name=gearman-statink-%(process_num)s
priority=1
numprocs=4
numprocs_start=2
autorestart=true
user=statink
